# config/supervisor/gsad.ini
[program:gsad]
command=/usr/local/bin/gsad.sh
user=root
directory=/var/run/
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/gsad.log