# config/supervisor/gvmd.ini
[program:gvmd]
command=/usr/sbin/gvmd --listen=0.0.0.0 --port=9390  -f
user=gvm
directory=/var/lib/gvm
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/gvmd.log