# config/supervisor/ospd-openvas.ini
[program:ospd-openvas]
command=/opt/atomicorp/bin/ospd-openvas -f --pid-file /var/run/ospd/ospd-openvas.pid --unix-socket=/var/run/ospd/ospd.sock --lock-file-dir /var/run/gvm/
user=gvm
directory=/var/lib/gvm
environment=HOME="/var/lib/gvm",USER="gvm",PYTHONPATH="/opt/atomicorp/lib/python3.6/site-packages"
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/ospd-openvas.log
