# config/supervisor/postgres.ini
[program:postgres]
command=/usr/local/bin/launch_postgres.sh postgres
user=postgres
directory=/var/lib/pgsql
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/postgres.log