# config/supervisor/redis.ini
[program:redis]
command=/usr/bin/redis-server /etc/redis.conf
user=gvm
directory=/var/lib/gvm
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/redis.log