# config/supervisor/startup.ini
[program:startup]
command=/usr/local/bin/startup.sh
user=root
autostart=true
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/supervisor/startup.log